# auth_dns/auth_dns.service
[Unit]
Description=Authoritative DNS Server
After=network.target
Wants=network.target

[Service]
Type=simple
# Starts as root to bind port 53, then drops to an unprivileged user in-process
# via -U (the binary keeps the bound socket fds across the drop).  To use it:
#   1) create the user:  useradd --system --no-create-home --shell /usr/sbin/nologin dnsd
#   2) make the log dir writable:  chown -R dnsd:dnsd @LOG_DIR@
#   3) append "-U dnsd" to ExecStart below.
# Alternatively run unprivileged from the start with systemd-native caps:
#   User=dnsd / Group=dnsd / AmbientCapabilities=CAP_NET_BIND_SERVICE
User=root
ExecStart=@INSTALL_BIN@/auth_dns -p 5335 -t 20 -u 127.0.0.1 -q 512
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3s
KillMode=mixed
TimeoutStopSec=15s
LimitNOFILE=65535

# --- Sandboxing / hardening ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectClock=true
ProtectHostname=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
# CAP_NET_BIND_SERVICE: bind port 53.  CAP_SETUID/SETGID: the in-process -U drop.
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_SETUID CAP_SETGID
ReadWritePaths=@LOG_DIR@ /run
ReadOnlyPaths=@SERVER_PATH@

[Install]
WantedBy=multi-user.target
